# cron/cron.hourly
#!/bin/bash
curl -I http://localhost:80/api/refresh/feed